# crontab
*/30 * * * * /var/www/html/script_shell.sh >> /var/www/html/log.txt
